# ops/systemd/pm-heartbeat.service
[Unit]
Description=Polymarket heartbeat guard
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/polymarket-bot/current
ExecStart=/opt/polymarket-bot/current/.venv/bin/python -m src.app run-heartbeat
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

