# ops/systemd/pm-quoteable-window-monitor.service
[Unit]
Description=Polymarket quoteable window monitor
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/polymarket-bot/current
ExecStart=/opt/polymarket-bot/current/.venv/bin/python -m src.app run-quoteable-window-monitor --selection-mode broad --iterations 0 --sleep-seconds 300 --cycle-minutes 5
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
