# ops/systemd/pm-strategy-mm.service
[Unit]
Description=Polymarket MM strategy worker
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/polymarket-bot/current
ExecStart=/opt/polymarket-bot/current/.venv/bin/python -m src.app run-strategy-mm
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

